# scripts/azure-dev/templates/quadlet/krav-kong.container
[Unit]
Description=Kravhantering Kong HSA proxy
Requires=krav-support-network.service krav-hsa-person-lookup-adapter.service
After=krav-support-network.service krav-hsa-person-lookup-adapter.service

[Container]
ContainerName=kong
Image=docker.io/kong/kong-gateway:3.15.0.0-20260702-ubuntu@sha256:f0516241dd354e66a9a8b4151c884bf40df8f8d850cec892adb4446cd210491f
Network=krav-support.network
PublishPort=127.0.0.1:18000:8000
Volume=/workspace/containers/kong/kong.yml:/kong/declarative/kong.yml:ro
Environment=KONG_ADMIN_ACCESS_LOG=/dev/stdout
Environment=KONG_ADMIN_ERROR_LOG=/dev/stderr
Environment=KONG_ADMIN_LISTEN=0.0.0.0:8001
Environment=KONG_DATABASE=off
Environment=KONG_DECLARATIVE_CONFIG=/kong/declarative/kong.yml
Environment=KONG_PREFIX=/tmp/kong
Environment=KONG_PROXY_ACCESS_LOG=/dev/stdout
Environment=KONG_PROXY_ERROR_LOG=/dev/stderr
Environment=KONG_PROXY_LISTEN=0.0.0.0:8000
HealthCmd=kong health
HealthInterval=10s
HealthRetries=12
HealthStartPeriod=10s
HealthTimeout=5s
HealthOnFailure=kill

[Service]
Environment=HOME=%h
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_HOME=%h/.local/share
Environment=XDG_RUNTIME_DIR=%t
Environment=CONTAINERS_CONF=%h/.config/containers/containers.conf
Environment=CONTAINERS_STORAGE_CONF=%h/.config/containers/storage.conf
Restart=on-failure
TimeoutStartSec=300

[Install]
WantedBy=default.target
