# containers/production/systemd/kravhantering-compose.service
# cSpell:words oneshot

[Unit]
Description=Kravhantering Podman Compose app node
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/kravhantering/current
Environment=PODMAN_COMPOSE_PROVIDER=podman-compose
ExecStart=/usr/bin/podman compose --env-file /etc/kravhantering/release.env -f compose/app-node-tls.compose.yml up -d
ExecStop=/usr/bin/podman compose --env-file /etc/kravhantering/release.env -f compose/app-node-tls.compose.yml down
TimeoutStartSec=180
TimeoutStopSec=90

[Install]
WantedBy=default.target
